# system/etc/profile.d/stellar-neofetch.sh
alias neofetch='/usr/bin/fastfetch --color $(/usr/libexec/fancy-fastfetch) -c /usr/share/ublue-os/stellarite/fastfetch.jsonc'
alias neowofetch='/usr/bin/fastfetch --color $(/usr/libexec/fancy-fastfetch) -c /usr/share/ublue-os/stellarite/fastfetch.jsonc'
alias hyfetch='/usr/bin/fastfetch --color $(/usr/libexec/fancy-fastfetch) -c /usr/share/ublue-os/stellarite/fastfetch.jsonc'
alias fastfetch='/usr/bin/fastfetch --color $(/usr/libexec/fancy-fastfetch) -c /usr/share/ublue-os/stellarite/fastfetch.jsonc'
